# deploy/systemd/mycelium-child-agent.service
[Unit]
Description=Mycelium Child Agent
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Mycelium
ExecStart=%h/Mycelium/scripts/run_child.sh
Restart=always
RestartSec=5
EnvironmentFile=-%h/.config/mycelium/child-agent.env

[Install]
WantedBy=default.target
